# service/log/run
#!/bin/sh
exec 2>&1
# documentation on how to use multilog https://manpages.debian.org/stretch/daemontools/multilog.8.en.html
